# rojo_scripts/rojo.crontab
# Minute   Hour   Day of Month       Month          Day of Week        Command    
# (0-59)  (0-23)     (1-31)    (1-12 or Jan-Dec)  (0-6 or Sun-Sat)                
    0        2          *              *                *            python3 /home/charles/rojo_scripts/pull_charlesreid1_data.py
